# files/cobbler-reposync.cron
#!/bin/bash

# Dynamically generated by Puppet
# Local modifications will be over-written

/usr/bin/cobbler reposync --tries=3 --no-fail &>> /dev/null
/usr/bin/cobbler hardlink &>> /dev/null
true
